# tools/email/systemd/protonmail-sync.service
[Unit]
Description=ProtonMail Sync Service
Documentation=man:protonmail(1)
After=network-online.target
Wants=network-online.target
# Only run if ProtonMail Bridge is running
ConditionPathExists=/usr/bin/protonmail-bridge

[Service]
Type=oneshot
Environment="PATH=%h/.local/bin:/usr/local/bin:/usr/bin:/bin"
Environment="DDATA=/data"
EnvironmentFile=%h/.config/protonmail/env

# Check if bridge is running before syncing
ExecStartPre=/bin/bash -c 'pgrep -f "protonmail-bridge|proton-bridge" > /dev/null || exit 1'

# Run the sync
ExecStart=%h/.local/bin/protonmail sync

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=protonmail-sync

# Security settings
NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/data/messages/email/protonmail

# Resource limits
MemoryMax=512M
TasksMax=50

[Install]
WantedBy=default.target
