# tools/email/systemd/protonmail-sync.timer
[Unit]
Description=ProtonMail Sync Timer (every 5 minutes)
Documentation=man:protonmail(1)
Requires=protonmail-sync.service

[Timer]
# Run every 5 minutes
OnBootSec=1min
OnUnitActiveSec=5min

# Allow timer to wake system from suspend (optional)
WakeSystem=false

# If the system was offline, don't run all missed timers at once
Persistent=false

[Install]
WantedBy=timers.target
